# sd-image/seed/fs/etc/systemd/system/network.target.wants/wpa_supplicant-configure.service
[Unit]
Description=horizon wpa supplicant configure service
Requires=sysinit.target
After=local-fs.target
Before=network-pre.target
RequiresMountsFor=/boot/firmware

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wpa-supplicant-configure
